# webserver.service
[Unit]
Description=Streaming Web Server

[Service]
ExecStart=/home/pi/webserver/webserver.sh
Restart=always

[Install]
WantedBy=default.target
